# templates/cron/launchd/com.brain.pattern-check.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0"><dict>
    <key>Label</key><string>com.brain.pattern-check</string>
    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>{{VAULT_PATH}}/06-Agent/cron/jobs/pattern-check.sh</string>
    </array>
    <key>WorkingDirectory</key><string>{{VAULT_PATH}}</string>
    <key>StartCalendarInterval</key>
    <dict>
        <key>Hour</key><integer>22</integer>
        <key>Minute</key><integer>30</integer>
    </dict>
    <key>StandardOutPath</key><string>{{VAULT_PATH}}/06-Agent/cron/logs/pattern-check.log</string>
    <key>StandardErrorPath</key><string>{{VAULT_PATH}}/06-Agent/cron/logs/pattern-check-error.log</string>
    <key>RunAtLoad</key><false/>
</dict></plist>
